# supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /home/docker/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx

[program:postgresql]
user=postgres
command=/usr/lib/postgresql/9.3/bin/postgres -D /var/lib/postgresql/9.3/main -c config_file=/etc/postgresql/9.3/main/postgresql.conf
process_name=%(program_name)s
stopsignal=INT
autostart=true
autorestart=true
